# agents/osx/com.memorable.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.memorable.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>AGENT_PATH/memorable-agent.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MEMORABLE_API</key>
        <string>http://memorable-alb-1679440696.us-west-2.elb.amazonaws.com</string>
        <key>MEMORABLE_USER</key>
        <string>alan</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/memorable-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/memorable-agent.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
